# deploy/webrtc-bridge.service
[Unit]
Description=WebRTC-to-RTSP Screen Share Bridge for Kodi
After=kodi.service network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/storage/kodi-screenshare/webrtc-bridge \
    -listen-addr :443 \
    -mediamtx-path ./mediamtx \
    -kodi-endpoint http://127.0.0.1:8080/jsonrpc \
    -kodi-username kodi \
    -kodi-password %KODI_PASSWORD% \
    -stream-host %STREAM_HOST%
WorkingDirectory=/storage/kodi-screenshare
Restart=on-failure
RestartSec=5

[Install]
WantedBy=kodi.target
